# files/services/dnstap/run
#!/usr/bin/env sh
export GOPATH
if [ "$RUN_DNSTAP" = true ]; then
  sleep 5
  exec /usr/sbin/gosu named:named bash -c "${GOPATH}/bin/dnstap -u /chroot/var/run/named/dnstap.sock"
fi
