# files/services/dnstap/run
#!/usr/bin/env sh
export GOPATH

if "$RUN_DNSTAP" ; then
  sleep 5
  exec su-exec named ${GOPATH}/bin/dnstap -u /chroot/var/run/named/dnstap.sock
fi
