# files/services/dnstap/run
#!/usr/bin/env sh
sleep 5
export GOPATH
exec /usr/sbin/gosu named:named bash -c "${GOPATH}/bin/dnstap -u /chroot/var/run/named/dnstap.sock"
