# files/services/dnstap/run
#!/usr/bin/env sh
export GOPATH
exec su-exec named ${GOPATH}/bin/dnstap -u /chroot/var/run/named/dnstap.sock
